# plasma-runner-js-eval.service
[Unit]
Description=KRunner plugin for JavaScript evaluation
After=graphical.target

[Install]
WantedBy=graphical.target

[Service]
ExecStart=node /path/to/program/lib/index.js
