# plasma-runner-js-eval.service
[Unit]
Description=KRunner JavaScript evaluation plugin
After=graphical.target

[Install]
WantedBy=graphical.target

[Service]
ExecStart=node /path/to/program/index.js
